# configs/supervisor/supervisor-app.conf
[program:uwsgi-app]
command=/usr/local/bin/uwsgi --ini /configs/uwsgi/uwsgi.ini
